# data-pipeline/fluid-indexer.service
[Unit]
Description=Fluid Protocol Real-Time Indexer
After=network.target clickhouse-server.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/RLD/data-pipeline
EnvironmentFile=-/home/ubuntu/RLD/data-pipeline/.env
ExecStart=/home/ubuntu/RLD/data-pipeline/.venv/bin/python3 /home/ubuntu/RLD/data-pipeline/scripts/fluid_realtime_indexer.py
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
